# contrib/services/harvester/run
#!/bin/bash
set -ex

/usr/bin/supervisord

echo "Looking for log"
while [ ! -f /var/log/syslog ]; do
    echo "."
    sleep 2
done
echo "Found the file"
tail -F /var/log/syslog
